# supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=npm start
redirect_stdout=true
